# supervisor/conf.d/import-data.conf
[program:import-data]
command=/root/install-scripts/import-data.sh
auto_start=true
autorestart=false
