# supervisor/conf.d/import-data.conf
[program:import-data]
command=bash -c "source /etc/profile.d/lizards-env.sh && /root/install-scripts/import-data.sh"
auto_start=true
autorestart=false
