# supervisor/conf.d/nginx.conf
[program:nginx]
command=/usr/sbin/nginx -g 'daemon off; master_process on;'
autostart=true
autorestart=true
